# root/etc/s6-overlay/s6-rc.d/init-user/run
#!/command/with-contenv sh

set -eu

PUID=${PUID:-999}
PGID=${PGID:-999}

groupmod --non-unique --gid "$PGID" steam
usermod --non-unique --uid "$PUID" steam

chown --recursive "$PUID:$PGID" ~steam

